# services/mlflow/crontab
* * * * * /usr/local/bin/python3 /app/hello.py >> /var/log/cron.log 2>&1 
